# server/upstart.conf
#!upstart
description "opendataheatmap"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec nodejs /root/projects/opendata-heatmap/server/server.js >> /var/log/opendataheatmap.log 2>&1
